# src/nodot/scripts/new_computer/at_home/backup/crontab
0 0 * * *   /home/backup/bin/do_backups.bash
